# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:backend]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=USE_VERTEX_AI="%(ENV_USE_VERTEX_AI)s",GOOGLE_CLOUD_PROJECT="%(ENV_GOOGLE_CLOUD_PROJECT)s",GOOGLE_CLOUD_LOCATION="%(ENV_GOOGLE_CLOUD_LOCATION)s"

[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="%(ENV_PORT)s",HOSTNAME="0.0.0.0",BACKEND_URL="http://localhost:8000"
